# deploy/canvas-mcp.service
[Unit]
Description=Canvas MCP Server (HTTP transport)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/canvas-mcp
ExecStart=/opt/canvas-mcp/.venv/bin/canvas-mcp-server --transport streamable-http --port 8819
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
